# system/crontab/crontab.as_a_ref
@reboot sleep 60 && nohup /home/nextcloud/distribusi-watcher.sh >> /home/nextcloud/distribusi-watcher-sh.log 2>&1

